# rules.d/50-radeon.rules
# /etc/udev/rules.d/50-radeon.rules
# vim: ft=udevrules

KERNEL=="card1", SUBSYSTEM=="drm", RUN+="/usr/bin/radeon"

