# rules.d/50-radeon.rules
# /etc/udev/rules.d/40-radeon.rules
# vim: ft=udevrules
KERNEL=="radeon", RUN+="/usr/bin/radeon"

